# systemd/renew.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/home/user/do-ssl/renew.sh
