# systemd/renew.service
[Unit]
Description=Let's Encrypt renewal

[Service]
Type=oneshot
ExecStart=/home/user/digitalocean-ssl/renew.sh
